# debian/wb-homa-modbus.postinst
#!/bin/sh

if [ -x "/etc/init.d/wb-homa-modbus" ] || [ -e "/etc/init/wb-homa-modbus.conf" ]; then
	if [ ! -e "/etc/init/wb-homa-modbus.conf" ]; then
		update-rc.d wb-homa-modbus defaults >/dev/null
	fi
	invoke-rc.d wb-homa-modbus start || true
fi

exit 0
# ignore debhelper

#DEBHELPER#
